# build/supervisor/server/httpsproxy.ini
[program:httpsproxy]
user=root
environment=SAAS_COMMON_CONFIG=test
command=/data/sa/service/httpsproxy/httpsproxy &
directory=/data/sa/service/httpsproxy

stopsignal=INT
startsecs=3
startretries=300
stopasgroup=true
killasgroup=true

autostart=true
autorestart=true
stdout_logfile=/data/supervisor/log/supervisor.log
stderr_logfile=/data/supervisor/log/supervisor_err.log
[supervisord]